# scripts/systemd/research-radar-daily.service
[Unit]
Description=Research Radar daily ingest, digest, and email

[Service]
Type=oneshot
WorkingDirectory=%h/Desktop/codes/research-scout
EnvironmentFile=%h/Desktop/codes/research-scout/.env
ExecStart=%h/Desktop/codes/research-scout/env/bin/research-radar run-all
