# scripts/systemd/research-radar-weekly.timer
[Unit]
Description=Run Research Radar weekly discovery on Monday at 06:30

[Timer]
OnCalendar=Mon *-*-* 06:30:00
Persistent=true

[Install]
WantedBy=timers.target
