# scripts/systemd/research-radar-weekly.service
[Unit]
Description=Research Radar weekly discovery

[Service]
Type=oneshot
WorkingDirectory=%h/Desktop/codes/research-scout
EnvironmentFile=%h/Desktop/codes/research-scout/.env
ExecStart=%h/Desktop/codes/research-scout/env/bin/research-radar discover-weekly
